# launchd_proposed/com.tailor.auto-flush-queue.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tailor.auto-flush-queue</string>
    <key>UserName</key>
    <string>jarvis</string>
    <key>GroupName</key>
    <string>staff</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/jarvis/tailor/scripts/services/run_auto_flush_queue.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/jarvis/tailor</string>
    <key>RunAtLoad</key>
    <false/>
    <key>KeepAlive</key>
    <false/>
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>14</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>20</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/jarvis/tailor/logs/auto_flush_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jarvis/tailor/logs/auto_flush_stderr.log</string>
</dict>
</plist>
