# launchd_proposed/com.tailor.secret-scan.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tailor.secret-scan</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/jarvis/tailor/scripts/maintenance/secret_scan.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/jarvis/tailor</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>TAILOR_HOME</key>
        <string>/Users/jarvis/tailor</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>0</integer>
        <key>Minute</key>
        <integer>15</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>

    <key>AbandonProcessGroup</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/jarvis/tailor/logs/secret_scan_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jarvis/tailor/logs/secret_scan_stderr.log</string>
</dict>
</plist>
